# home/Library/LaunchAgents/com.claude.schedule.ob-daily-log.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.ob-daily-log</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.claude/schedules/ob-daily-log/run-scheduler.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.claude/logs/ob-daily-log.out</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.claude/logs/ob-daily-log.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
