# home/Library/LaunchAgents/com.claude.schedule.tom-doerr-url-to-x-post.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.schedule.tom-doerr-url-to-x-post</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.claude/schedules/tom-doerr-url-to-x-post/run-scheduler.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.claude/logs/tom-doerr-url-to-x-post.out</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.claude/logs/tom-doerr-url-to-x-post.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
